# ansible/playbooks/templates/etc/init.d/thost.j2
#!/bin/bash
# chkconfig: 2345 20 80
# description: Tiddlyhost

START_DIR={{ home_dir }}

start() {
    docker compose up --detach
}

stop() {
    # The --volumes option is so that the rails-static volume is removed on
    # shutdown. It will be recreated again on startup which is what we want.
    docker compose down --volumes
}

# Recreate app, jobs, and web containers, leaving db and cache running.
# Web must be stopped too so the rails-static volume can be removed
# and recreated with fresh static assets from the new app image.
fast_restart() {
    docker compose stop app jobs web
    docker compose rm -f app jobs web
    docker volume rm th_rails_static 2>/dev/null || true
    docker compose up -d --no-deps app jobs web
}

cd $START_DIR

case "$1" in
    start)
        start
        ;;

    stop)
        stop
        ;;

    restart)
        stop
        start
        ;;

    fast-restart)
        fast_restart
        ;;

esac

exit 0
