# ansible/playbooks/templates/etc/systemd/system/thost.service.d/override.conf.j2
[Unit]
RequiresMountsFor={{ home_dir }}
